# etc/update-motd.d/98-footer
#!/bin/sh

/usr/games/cowsay "$(cat $(find  /etc/update-motd.d/quotes -type f | shuf -n 1))"
