# etc/systemd/system/reboot.timer
[Unit]
Description=Reboots system on schedule

[Timer]
OnCalendar=Sat
Persisten=true

[Install]
WantedBy=timers.target

