# etc/systemd/system/updates.service
[Unit]
Description=Update system

[Service]
Type=simple
ExecStart=/usr/local/sbin/upgrade_system.sh
Environment="HOME=/root"

