# etc/systemd/system/updates.timer
[Unit]
Description=Updates system on schedule

[Timer]
OnCalendar=daily
Persisten=true

[Install]
WantedBy=timers.target

